# teeunit/server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:teeworlds]
command=/opt/teeworlds/teeworlds_srv -f /opt/teeworlds/teeworlds.cfg
directory=/opt/teeworlds
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/teeworlds.log
stderr_logfile=/var/log/supervisor/teeworlds_err.log
priority=10

[program:fastapi]
command=uvicorn teeunit.server.app:app --host 0.0.0.0 --port %(ENV_PORT)s --log-level info
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fastapi.log
stderr_logfile=/var/log/supervisor/fastapi_err.log
priority=20
startsecs=5
